# rbn/rbn-filter.socket
[Socket]
ListenStream=7000
Accept=yes

[Install]
WantedBy=sockets.target
