# udev/rbn-filter@.service
[Service]
Type=exec
ExecStart=-/home/myon/shack/rbn/rbn-filter
#StandardInput=socket
StandardOutput=socket
