# udev/tlf-log.path
[Path]
PathChanged=/home/myon/log/tlf/contest.log

[Install]
WantedBy=default.target
